# systemd/pdrive.service
[Unit]
Description=Proton Drive sync daemon
Documentation=https://github.com/YOUR_USERNAME/proton-drive
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/pdrive-daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
